# systemd.mcman.target
[Unit]
Description=Minecraft Managed Instances
Wants=

[Install]
WantedBy=multi-user.target